# config/supervisor.conf
[program:gunicorn]
command=/home/khazidhea/.virtualenvs/loans/bin/gunicorn core.wsgi -c gunicorn.conf
directory=/srv/loans
autostart=true
autorestart=true
redirect_stderr=true

